# merox-agent.service
[Unit]
Description=Merox Infrastructure Agent
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
WorkingDirectory=/srv/merox-agent
ExecStart=/root/.local/share/mise/installs/python/3.14.3/bin/python3 /srv/merox-agent/service.py
EnvironmentFile=/srv/merox-agent/.env
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
