# merox-agent.service
[Unit]
Description=Merox Infrastructure Agent
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
WorkingDirectory=/srv/merox-agent
ExecStart=/srv/merox-agent/.venv/bin/python3 /srv/merox-agent/service.py
EnvironmentFile=/srv/merox-agent/.env
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
